# postgres/container_files/download-db.service
[Unit]
Description=Download IMDB

[Service]
ExecStart=/usr/local/bin/download-db.sh

[Install]
WantedBy=multi-user.target